# docker/laravel.crontab
* * * * * cd /var/www/v2board && php artisan schedule:run >> /dev/null 2>&1
